# rate-mirrors/rate-mirrors-apply.service
[Unit]
Description=Install the mirrorlist that was generated by rate-mirrors-gather.service and saved to /home/%USER%/.mirrorlist.ranked.
Requires=rate-mirrors-gather.service
After=rate-mirrors-gather.service

[Service]
ExecStart=bash -c "cp /home/%USER%/.mirrorlist.ranked /etc/pacman.d/mirrorlist && rm /home/%USER%/.mirrorlist.ranked"
Type=oneshot

